# image/downstage-provision.service
[Unit]
Description=Downstage OS first-boot provisioning
After=local-fs.target
Before=NetworkManager.service avahi-daemon.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/firstboot-provision.py
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
